# container/root/etc/services.d/scarlett_os/run
#!/usr/bin/execlineb -P

with-contenv jhbuild run python --version
