# container/root/etc/services.d/scarlett_os/run
#!/usr/bin/execlineb -P

with-contenv /home/pi/.virtualenvs/scarlett_os/bin/python3 --version
